# hoarder/rootfs/etc/s6-overlay/s6-rc.d/chromium/run
#!/command/with-contenv bashio

if [[ $CHROMIUM_ENABLED == "false" ]]; then
  bashio::log.info "Chromium is disabled"
else
  bashio::log.info "Starting Chromium..."
  chromium-browser --headless --remote-debugging-address=0.0.0.0 --remote-debugging-port=9222 $CHROMIUM_FLAGS
fi
