# traefik/rootfs/etc/s6-overlay/s6-rc.d/traefik/run
#!/command/with-contenv bashio

/traefik
