# actual/rootfs/etc/s6-overlay/s6-rc.d/actual/run
#!/command/with-contenv bashio

bashio::log.info "Starting Actual..."

cd /app
exec yarn start:server
