# wg-easy/rootfs/etc/s6-overlay/s6-rc.d/wg-easy/run
#!/command/with-contenv bashio

bashio::log.info "Starting WireGuard Easy..."
cd /app
node server/index.mjs
