# traefik/rootfs/etc/s6-overlay/s6-rc.d/init-traefik/run
#!/command/with-contenv bashio

mkdir -p /etc/traefik
printf -- "$(bashio::config '"traefik.yml"')" > /etc/traefik/traefik.yml
printf -- "$(bashio::config '"dynamic_conf.yml"')" > /etc/traefik/dynamic_conf.yml
