# wg-easy/rootfs/etc/s6-overlay/s6-rc.d/configure-wg-easy/run
#!/command/with-contenv bashio

while IFS= read -r key; do
  printf -- "$(bashio::config $key)" > "/var/run/s6/container_environment/$key"
done < <(jq -r 'keys[]' /data/options.json)
